# infra/bronze/roles/minio/templates/minio.service.j2
[Unit]
Description=MinIO
After=network-online.target

[Service]
Type=simple
User=minio-user
EnvironmentFile=-/etc/default/minio
ExecStart=/usr/local/bin/minio server {{ minio_data_dir }} \
  --address :{{ minio_port }} \
  --console-address :{{ minio_console_port }}
Restart=always

[Install]
WantedBy=multi-user.target
